# etc/rc.local
#Put your custom commands here that should be executed once
#the system init finished

# Essential scripts for robot

python /home/python-server/run.py &
#python /home/postman.py &
#python /home/receive.py &

# put your own scripts below this line

exit 0
