# etc/rc.local
#Put your custom commands here that should be executed once
#the system init finished

#Essencial script for robot

#python /home/startup.py
python /home/python-server/run.py &
#python /home/postman.py &
#python /home/receive.py &

#put your own scripts bellow this line

exit 0
